# .config/systemd/user/calcurse-notify.timer
[Unit]
Description=Timer for calcurse-notify.service

[Timer]
Persistent=true
OnBootSec=120
OnUnitActiveSec=7200

[Install]
WantedBy=timers.target
